# sms_sender/deployement/supervisor.ini
[program:sms_sender_service]
command=/home/sendersm/public_html/sms_sender/deployement/gunicorn_service
user=root
stdout_logfile=/home/sendersm/logs/gunicorn_supervisor.log
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8
